# ansible.cron
# Ansible pull cron job
*/10 * * * * root /usr/sbin/kamos-ansible-pull