# ansible.cron
#!/bin/bash

# Run ansible-pull with authentication for a private GitHub repository
/usr/bin/ansible-pull -U https://<username>:<personal_access_token>@github.com/<username>/<repository>.git /workspaces/kamos-ansible/ansible-repo/pull.yml